# containers/ddev-webserver/files/etc/supervisor/conf.d/mailhog.conf
[program:mailhog]
command=/usr/local/bin/mailhog
autostart=true
autorestart=true
startretries=10
directory=/var/www/html
redirect_stderr=true
