# containers/ddev-webserver/files/etc/supervisor/conf.d/mailhog.conf
[program:mailhog]
command=/usr/local/bin/mailhog
autorestart=true
startretries=10
