# containers/ddev-webserver/ddev-webserver-dev-base-files/etc/supervisor/conf.d/mailhog.conf
[program:mailpit]
command=/usr/local/bin/mailpit
autorestart=true
startretries=3
